# .config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper using swww every 12 hours

[Service]
ExecStart=%h/.local/bin/change-wallpaper
